# akiflow-sync/akiflow-sync.service
[Unit]
Description=Akiflow Sync Daemon
After=network.target
PartOf=nanoclaw.service

[Service]
Type=simple
WorkingDirectory=%h/code/yonibot/gabay
EnvironmentFile=%h/code/yonibot/gabay/.env
ExecStartPre=/usr/bin/npm run build --prefix %h/code/yonibot/gabay/akiflow-sync
ExecStart=/usr/bin/node %h/code/yonibot/gabay/akiflow-sync/dist/daemon.js
Restart=always
RestartSec=5
StandardOutput=append:%h/code/yonibot/gabay/akiflow/akiflow-sync.log
StandardError=append:%h/code/yonibot/gabay/akiflow/akiflow-sync.log

[Install]
WantedBy=default.target
